# wesai/open-api/task/docs/supervisor/local-open-api-task.ini
[group:local_open_api_task_single]
programs=NotifySoftykt

[program:NotifySoftykt]
command=/usr/bin/php /home/local/code/open-api/task/webroot/index.php softykt/NotifySoftykt
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-api/logs/task/notifySoftykt.log
stdout_logfile_maxbytes=20MB
