# wesai/open-operation/task/docs/supervisor/local-open-operation-task.ini
[group:local_open_operation_task]
programs=Subscribe,Unsubscribe

[program:Subscribe]
command=/usr/bin/php /home/local/code/open-operation/task/webroot/index.php wxhandler/event/subscribe
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-operation/logs/task/wxhandler_event_subscribe.log
stdout_logfile_maxbytes=10MB

[program:Unsubscribe]
command=/usr/bin/php /home/local/code/open-operation/task/webroot/index.php wxhandler/event/unsubscribe
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-operation/logs/task/wxhandler_event_unsubscribe.log
stdout_logfile_maxbytes=10MB

