# wesai/open-operation/task/docs/supervisor/local-open-operation-task.single.ini
[group:local_open_operation_task_single]
programs=CancelVoucher,CreateVoucher,OverdueVoucher

[program:CancelVoucher]
command=/usr/bin/php /home/local/code/open-operation/task/webroot/index.php voucher/CancelVoucher
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-operation/logs/task/voucher_CancelVoucher.log
stdout_logfile_maxbytes=20MB

[program:CreateVoucher]
command=/usr/bin/php /home/local/code/open-operation/task/webroot/index.php voucher/CreateVoucher
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-operation/logs/task/voucher_CreateVoucher.log
stdout_logfile_maxbytes=20MB

[program:OverdueVoucher]
command=/usr/bin/php /home/local/code/open-operation/task/webroot/index.php voucher/OverdueVoucher
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-operation/logs/task/voucher_OverdueVoucher.log
stdout_logfile_maxbytes=20MB