# wesai/open-book/task/docs/supervisor/local-open-book-task.ini
book[group:local_open_contest_task]
programs = order_CloseOrder,order_PayCompleted,order_PayResultCheck

[program:order_CloseOrder]
command = /usr/bin/php /home/local/code/open-book/task/webroot/index.php order/CloseOrder
redirect_stderr = true
numprocs = 3
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-book/logs/task/order_CloseOrder_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-book/logs/task/order_CloseOrder_stderr.log
autorestart = true


[program:order_PayCompleted]
command = /usr/bin/php /home/local/code/open-book/task/webroot/index.php order/PayCompleted
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-book/logs/task/order_PayCompleted_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-book/logs/task/order_PayCompleted_stderr.log
autorestart = true


[program:order_PayResultCheck]
command = /usr/bin/php /home/local/code/open-book/task/webroot/index.php order/PayResultCheck
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-book/logs/task/order_PayResultCheck_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-book/logs/task/order_PayResultCheck_stderr.log
autorestart = true

[program:venue_CreateTimes]
command = /usr/bin/php /home/local/code/open-book/task/webroot/index.php venue/CreateTimes
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-book/logs/task/order_PayResultCheck_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-book/logs/task/order_PayResultCheck_stderr.log
autorestart = true

