# wesai/open-book/task/docs/supervisor/local-open-book-task.single.ini
[group:local_open_contest_task_single]
programs = order_LoadOrderExpired

[program:order_LoadOrderExpired]
command = /usr/bin/php /home/local/code/open-book/task/webroot/index.php order/LoadOrderExpired
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-book/logs/task/order_LoadOrderExpired_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-book/logs/task/order_LoadOrderExpired_stderr.log
autorestart = true
