# wesai/open-little/task/docs/supervisor/local-open-little-task.ini
[group:local_open_little_task_single]
programs=Notify

[program:NotifySoftykt]
command=/usr/bin/php /home/local/code/open-little/task/webroot/index.php Notify/index
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-little/logs/task/notify.log
stdout_logfile_maxbytes=20MB
