# wesai/open-little/task/docs/crontab/open-little-crontab.single.config
# 每1个小时更新token
0 */1 * * * /bin/php /home/local/code/open-little/task/webroot/index.php GetToken/index