# wesai/open-contest/task/docs/supervisor/local-open-contest-task.single.ini
[group:local_open_contest_task_single]
programs = batch_order_LoadExpired,softykt_PlaceOrder

[program:batch_order_LoadExpired]
command = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php batch/order/LoadExpired
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-contest/logs/task/batch_order_LoadExpired_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-contest/logs/task/batch_order_LoadExpired_stderr.log
autorestart = true

[program:softykt_PlaceOrder]
command = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php softykt/PlaceOrder
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-contest/logs/task/softykt_PlaceOrder.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-contest/logs/task/softykt_PlaceOrder.log
autorestart = true
