# wesai/open-contest/task/docs/supervisor/local-open-contest-task.ini
[group:local_open_contest_task]
programs = batch_order_CheckPayResult,batch_order_Close,batch_order_Fail,batch_invite_code_Create


[program:batch_order_CheckPayResult]
command                 = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php batch/order/CheckPayResult
redirect_stderr         = true
numprocs                = 1
process_name            = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile          = /home/local/code/open-contest/logs/task/batch_order_CheckPayResult_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile          = /home/local/code/open-contest/logs/task/batch_order_CheckPayResult_stderr.log
autorestart             = true


[program:batch_order_Close]
command                 = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php batch/order/Close
redirect_stderr         = true
numprocs                = 3
process_name            = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile          = /home/local/code/open-contest/logs/task/batch_order_Close_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile          = /home/local/code/open-contest/logs/task/batch_order_Close_stderr.log
autorestart             = true


[program:batch_order_Fail]
command                 = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php batch/order/Fail
redirect_stderr         = true
numprocs                = 3
process_name            = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile          = /home/local/code/open-contest/logs/task/batch_order_Fail_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile          = /home/local/code/open-contest/logs/task/batch_order_Fail_stderr.log
autorestart             = true


[program:batch_invite_code_Create]
command                 = /usr/bin/php /home/local/code/open-contest/task/webroot/index.php batch/invite_code/Create
redirect_stderr         = true
numprocs                = 1
process_name            = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile          = /home/local/code/open-contest/logs/task/batch_invite_code_Create_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile          = /home/local/code/open-contest/logs/task/batch_invite_code_Create_stderr.log
autorestart             = true

