# wesai/open-contest/task/docs/crontab/open-contest-crontab.single.config

# 每天凌晨3:00发布统计前一天数据事件
0 3 * * * /bin/php /home/local/code/open-contest/task/webroot/index.php batch/analysis/LoadAnalysisTarget
