# wesai/open-weixin/task/docs/supervisor/local-open-weixin-task-qywx.ini
[group:local_open_weixin_task_qywx]
programs = qywx_auth_LoadAuthNeedRefreshToken,qywx_auth_RefreshToken,qywx_corp_LoadCorpNeedRefreshTicket,qywx_corp_RefreshTicket,qywx_provider_RefreshToken,qywx_suite_LoadSuiteNeedRefreshToken,qywx_suite_RefreshToken,qywx_suite_ConsumeCallbackEvent

[program:qywx_auth_LoadAuthNeedRefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/auth/LoadAuthNeedRefreshToken
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_LoadAuthNeedRefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_LoadAuthNeedRefreshToken_stderr.log
autorestart = true


[program:qywx_auth_RefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/auth/RefreshToken
redirect_stderr = true
numprocs = 3
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_RefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_RefreshToken_stderr.log
autorestart = true


[program:qywx_corp_LoadCorpNeedRefreshTicket]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/corp/LoadCorpNeedRefreshTicket
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_LoadCorpNeedRefreshTicket_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_LoadCorpNeedRefreshTicket_stderr.log
autorestart = true


[program:qywx_corp_RefreshTicket]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/corp/RefreshTicket
redirect_stderr = true
numprocs = 3
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_RefreshTicket_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_RefreshTicket_stderr.log
autorestart = true


[program:qywx_provider_RefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/provider/RefreshToken
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_provider_RefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_provider_RefreshToken_stderr.log
autorestart = true


[program:qywx_suite_LoadSuiteNeedRefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/suite/LoadSuiteNeedRefreshToken
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_LoadSuiteNeedRefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_LoadSuiteNeedRefreshToken_stderr.log
autorestart = true


[program:qywx_suite_RefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/suite/RefreshToken
redirect_stderr = true
numprocs = 3
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_RefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_RefreshToken_stderr.log
autorestart = true


[program:qywx_suite_ConsumeCallbackEvent]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/suite/ConsumeCallbackEvent
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_ConsumeCallbackEvent_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_ConsumeCallbackEvent_stderr.log
autorestart = true

