# wesai/open-weixin/task/docs/supervisor/local-open-weixin-task-qywx.single.ini
[group:local_open_weixin_task_qywx_single]
programs = qywx_auth_LoadAuthNeedRefreshToken,qywx_corp_LoadCorpNeedRefreshTicket,qywx_suite_LoadSuiteNeedRefreshToken

[program:qywx_auth_LoadAuthNeedRefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/auth/LoadAuthNeedRefreshToken
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_LoadAuthNeedRefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_auth_LoadAuthNeedRefreshToken_stderr.log
autorestart = true

[program:qywx_corp_LoadCorpNeedRefreshTicket]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/corp/LoadCorpNeedRefreshTicket
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_LoadCorpNeedRefreshTicket_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_corp_LoadCorpNeedRefreshTicket_stderr.log
autorestart = true

[program:qywx_suite_LoadSuiteNeedRefreshToken]
command = /usr/bin/php /home/local/code/open-weixin/task/webroot/index.php qywx/suite/LoadSuiteNeedRefreshToken
redirect_stderr = true
numprocs = 1
process_name = %(program_name)s_%(process_num)d
stdout_logfile_maxbytes = 20MB
stdout_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_LoadSuiteNeedRefreshToken_stdout.log
stderr_logfile_maxbytes = 20MB
stderr_logfile = /home/local/code/open-weixin/logs/task/qywx_suite_LoadSuiteNeedRefreshToken_stderr.log
autorestart = true
