# wesai/open-weixin/task/docs/supervisor/local-open-weixin-task-component.single.ini
[group:local_open_weixin_task_single]
programs=AuthorizerWillRefresh

[program:AuthorizerWillRefresh]
command=/bin/php /home/local/code/open-weixin/task/webroot/index.php refresh/componentAuthorizerWillRefresh
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/refresh_AuthorizerWillRefresh.log
stdout_logfile_maxbytes=10MB

