# wesai/open-weixin/task/docs/supervisor/local-open-weixin-task-component.ini
[group:local_open_weixin_task]
programs=AuthorizerRefreshToken,EventSubscribe,EventUnsubscribe,MsgText,Card_not_pass_check,Card_pass_check,Card_pay_order,Card_sku_remind,Submit_membercard_user_info,Update_member_card,User_consume_card,User_del_card,User_enter_session_from_card,User_get_card,User_gifting_card,User_pay_from_pay_cell,User_view_card

[program:AuthorizerRefreshToken]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php refresh/componentAuthorizerRefreshToken
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/refresh_AuthorizerRefreshToken.log
stdout_logfile_maxbytes=10MB

[program:EventSubscribe]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/subscribe
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_subscribe.log
stdout_logfile_maxbytes=10MB

[program:EventUnsubscribe]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/unsubscribe
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_unsubscribe.log
stdout_logfile_maxbytes=10MB

[program:MsgText]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/msg/text
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_msg_text.log
stdout_logfile_maxbytes=10MB

[program:Card_not_pass_check]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/card_not_pass_check
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_card_not_pass_check.log
stdout_logfile_maxbytes=10MB

[program:Card_pass_check]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/card_pass_check
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_card_pass_check.log
stdout_logfile_maxbytes=10MB

[program:Card_pay_order]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/card_pay_order
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_card_pay_order.log
stdout_logfile_maxbytes=10MB

[program:Card_sku_remind]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/card_sku_remind
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_card_sku_remind.log
stdout_logfile_maxbytes=10MB

[program:Submit_membercard_user_info]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/submit_membercard_user_info
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_submit_membercard_user_info.log
stdout_logfile_maxbytes=10MB

[program:Update_member_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/update_member_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_update_member_card.log
stdout_logfile_maxbytes=10MB

[program:User_consume_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_consume_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_user_consume_card.log
stdout_logfile_maxbytes=10MB

[program:User_del_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_del_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_cuser_del_card.log
stdout_logfile_maxbytes=10MB

[program:User_enter_session_from_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_enter_session_from_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_card_not_pass_check.log
stdout_logfile_maxbytes=10MB

[program:User_get_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_get_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_user_get_card.log
stdout_logfile_maxbytes=10MB

[program:User_gifting_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_gifting_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_user_gifting_card.log
stdout_logfile_maxbytes=10MB

[program:User_pay_from_pay_cell]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_pay_from_pay_cell
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_user_pay_from_pay_cell.log
stdout_logfile_maxbytes=10MB

[program:User_view_card]
command=/usr/bin/php /home/local/code/open-weixin/task/webroot/index.php wxhandler/event/user_view_card
process_name=%(program_name)s_%(process_num)s
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/home/local/code/open-weixin/logs/task/wxhandler_event_user_view_card.log
stdout_logfile_maxbytes=10MB

