# wesai/open-book/task/docs/crontab/open-book-crontab.config
#每分钟重建一次索引
*/1 * * * * /usr/bin/indexer -c /home/local/code/open-book/task/application/controllers/sphinx/conf/sphinx.conf --all --rotate

# 每天凌晨3:00重启searchd服务
0 3 * * * /usr/bin/bash /home/local/code/open-book/task/application/controllers/sphinx/restart.sh