# wesai/open-book/task/docs/crontab/open-book-crontab.single.config
# 每天凌晨0:00生成场地的res
0 0 * * * /bin/php /home/local/code/open-book/task/webroot/index.php venue/LoadRuleRes