# wesai/open-api/task/docs/crontab/open-api-crontab.single.config
# 每1个小时更新token
0 */1 * * * /bin/php /home/local/code/open-api/task/webroot/index.php softykt/GetToken/index