# wesai/open-weixin/task/docs/crontab/open-weixin-crontab.config
0 */1 * * * /bin/php /home/local/code/open-weixin/task/webroot/index.php gen/componentToken >> /home/local/code/open-weixin/logs/task/crontab.gen.componentToken 2>&1
